# root/etc/services.d/flexget/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc /usr/local/bin/flexget -c /config/config.yml --loglevel ${FLEXGET_LOG_LEVEL} daemon start --autoreload-config

